# systemd/rag-api.service
[Unit]
Description=RAG API Server (FastAPI)
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=rag
Group=rag
WorkingDirectory=/opt/rag-system

Environment="PATH=/opt/rag-system/venv/bin"
EnvironmentFile=/opt/rag-system/.env

ExecStart=/opt/rag-system/venv/bin/uvicorn src.api.main:app \
  --host 0.0.0.0 \
  --port 8001 \
  --workers 1 \
  --log-level info

Restart=always
RestartSec=10

LimitNOFILE=65535

[Install]
WantedBy=multi-user.target